# packages/DEBIAN/prerm
#!/bin/sh
set -e

echo "Pre-Removal Macro"

systemctl disable __PACKAGE_NAME__ || true
systemctl stop __PACKAGE_NAME__ || true

echo "Pre-Installation Macro called with argument: $1 ." >&2

#rm -f __PACKAGE_DIR__/config/settings.cfg

exit 0
